# linux-live/scripts/00-minios/rootcopy-install/usr/lib/systemd/system/polkit-actions-monitor.service
[Unit]
Description=React to modifications in /usr/share/polkit-1/actions

[Service]
ExecStart=/usr/bin/polkit-actions-handler
