# linux-live/scripts/00-minios/rootcopy-install/usr/lib/systemd/system/polkit-actions-monitor.path
[Unit]
Description=Monitor /usr/share/polkit-1/actions for changes

[Path]
PathModified=/usr/share/polkit-1/actions

[Install]
WantedBy=multi-user.target
